# playbooks/roles/jenkins_agent/templates/jenkins-agent.service.j2
[Unit]
Description=Jenkins Agent worker
After=syslog.target network.target

[Service]
Type=simple
User=jenkins
Group=jenkins
WorkingDirectory=/home/jenkins
Restart=always
RestartSec=10
Environment="JAVA_HOME={{ java_home }}"
Environment="CONFIG_URL={{ config_url }}"
ExecStart={{ script_path }}

[Install]
WantedBy=multi-user.target
